# ros2_ws/src/perception_costmap/deploy/percept-stack.service
# Boot autostart for the whole perception/camera/streaming stack.
# Install (as root on dinosaur):
#   cp percept-stack.service /etc/systemd/system/
#   systemctl daemon-reload && systemctl enable --now percept-stack
# Stop everything: systemctl stop percept-stack
[Unit]
Description=AVL perception stack (tmux 'percept': sensors, ZEDs, costmap, EKF, viz, streams)
After=network-online.target zed_x_daemon.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=dinosaur
WorkingDirectory=/home/dinosaur
ExecStart=/bin/bash /home/dinosaur/carla-nav2-avl/ros2_ws/src/perception_costmap/deploy/boot_stack.sh
ExecStop=/usr/bin/tmux -L percept kill-session -t percept
TimeoutStartSec=420

[Install]
WantedBy=multi-user.target
